# scripts/com.fsad.playbook-updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fsad.playbook-updater</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/theobeack/Desktop/AI/fsad_playbook/scripts/run-playbook-updater.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/theobeack/Desktop/AI/fsad_playbook</string>

    <!-- Run daily at 08:00 local time. Adjust Hour/Minute as needed. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- If the Mac was asleep at the scheduled time, run at next wake. -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/theobeack/Desktop/AI/fsad_playbook/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/theobeack/Desktop/AI/fsad_playbook/logs/launchd.err.log</string>
</dict>
</plist>
